# example/config/supervisor_app.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /example/config/uwsgi.ini
